# other/95-monitors.rules

ACTION=="change", RUN+="/home/vince/bin/disp-auto.sh"
